# src/src/Rc.DiscordBot.Worker.Linux/discordMusicBot.service
[Unit]
Description=Discord Music Bot

[Service]
Type=notify
ExecStart=/usr/sbin/DiscordMusicBot/DiscordMusicBot

[Install]
WantedBy=multi-user.target